# services/glusterd/run
#!/usr/bin/env bash

/usr/sbin/glusterd -p /var/run/glusterd.pid
