# services/glusterd/run
#!/usr/bin/env bash
service glusterfs-server start
